# conf/esmond-mesh-2-tsds.cron
*/15 * * * * netsage /usr/bin/perl /usr/bin/esmond2tsds -c /etc/netsage/esmond-mesh-2-tsds/config.xml >/dev/null 2>&1
